# root/usr/local/lib/systemd/system/containerd.service.d/conditional.conf
# Only enable containerd if it is the selected container engine.
[Service]
ExecCondition=test ${CONTAINER_ENGINE} = containerd
