# root/usr/local/lib/systemd/system/containerd.service.d/target.conf
# Overrride containerd to install to rancher-desktop.target
[Install]
WantedBy=
WantedBy=rancher-desktop.target

[Unit]
StopPropagatedFrom=rancher-desktop.target
