# root/usr/local/lib/systemd/system/lima-init.service.d/requires-cidata.conf
# On non-WSL platforms, we require /mnt/lima-cidata (for lima.env) before
# running lima-init.service

[Unit]
RequiresMountsFor=/mnt/lima-cidata
