# root/usr/local/lib/systemd/system/lima-init.service.d/requires-rd-init.conf
# Lima-init.service requires rd-init.service, except on WSL where rd-init.service
# is unneeded.

[Unit]
Wants=rd-init.service
After=rd-init.service
