# root/usr/local/lib/systemd/system/k3s.service.d/command-override.conf
# Override the k3s command line
[Service]
ExecStart=
# This should be kept in sync with `network-namespace.conf`.
ExecStart=/usr/local/bin/k3s server \
    --https-listen-port $KUBERNETES_PORT \
    $K3S_CONTAINER_ENGINE_OPTIONS
