# root/usr/local/lib/systemd/system/k3s.service.d/conditional.conf
# Only enable k3s if desired
[Service]
ExecCondition=test ${KUBERNETES_ENABLED} = true
