# root/usr/local/lib/systemd/system/k3s.service.d/dependencies.conf
# Makes sure k3s.service is started after it is installed and the container engine is ready.
[Unit]
Requires=k3s-install.service container-engine.target
After=k3s-install.service container-engine.target
