# root/usr/local/lib/systemd/system/k3s.service.d/environment.conf
# Import Rancher Desktop configuration environment variables
[Unit]
Requires=params.service
After=params.service
[Service]
EnvironmentFile=-/etc/sysconfig/rancher-desktop
