# root/usr/local/lib/systemd/system/k3s.service.d/network-namespace.conf
# Make k3s join the network namespace (for WSL).
[Unit]
Requires=setup-network.service
JoinsNamespaceOf=setup-network.service
