# root/usr/local/lib/systemd/system/k3s.service.d/target.conf
# We want to override WantedBy= to rancher-desktop.target, but that cannot be
# done in an override file - https://github.com/systemd/systemd/issues/33437
# Instead, we just set StopPropagatedFrom instead.
[Unit]
StopPropagatedFrom=rancher-desktop.target
