# root/usr/local/lib/systemd/system/local.service
[Unit]
Description=Rancher Desktop Provisioning Scripts

[Service]
Type=oneshot
ExecStart=/bin/true

[Install]
WantedBy=rancher-destkop.target
