# root/usr/local/lib/systemd/system/docker.service.d/allowed-images.conf
# Set HTTPS_PROXY if /usr/local/openresty/nginx/conf/allowed-images.conf exists.

[Service]
EnvironmentFile=-/etc/sysconfig/allowed-images.env
ExecStartPre=/bin/sh -c '\
    set -x ; \
    if $ALLOWED_IMAGES ; then \
        echo HTTPS_PROXY=http://127.0.0.1:3128 > /etc/sysconfig/allowed-images.env ; \
    else \
        rm -f /etc/sysconfig/allowed-images.env ; \
    fi \
    '
