# root/usr/local/lib/systemd/system/docker.service.d/conditional.conf
# Only enable docker if desired
[Service]
ExecCondition=test ${CONTAINER_ENGINE} = moby
