# root/usr/local/lib/systemd/system/docker.service.d/permissions.conf
# Override docker socket permissions so the user is able to list containers.
# Note that dockerd overwrites permissions, so we can't just set this as an
# override of docker.socket.

[Service]
ExecStartPost=chmod a+rw /var/run/docker.sock
