# root/usr/local/lib/systemd/system/openresty.service.d/ca-certs.conf
# Compatibility with alpine's ssl cert locations
[Service]
# This is an ExecCondition to start before ExecStartPre
ExecCondition=ln -sf /etc/ssl/ca-bundle.pem /etc/ssl/certs/ca-certificates.crt
