# root/usr/local/lib/systemd/system/cri-dockerd.service.d/network-namespace.conf
[Unit]
Requires=network-namespace.service network-setup.service
After=network-setup.service

[Service]
ExecStart=
ExecStart=/usr/local/bin/wsl-exec \
  /usr/local/bin/cri-dockerd \
  --container-runtime-endpoint unix:///run/cri-dockerd.sock \
  --network-plugin=cni \
  --cni-bin-dir=/usr/libexec/cni \
  --cni-conf-dir=/etc/cni/net.d \
  --cni-cache-dir=/var/lib/cni/cache
