# root/usr/local/lib/systemd/system/console-getty.service.d/exclude-containers.conf
# console-getty.service does not start correctly when running in a container
# (because they don't have separate consoles); disable in that case.
[Unit]
ConditionVirtualization=!container
