# root/usr/local/lib/systemd/system/local.service.d/network-namespace.conf
[Unit]
Requires=network-namespace.service network-setup.service
After=network-setup.service

[Service]
ExecStart=
ExecStart=/usr/local/bin/wsl-exec /bin/true
