# root/usr/local/lib/systemd/system/buildkitd.service.d/conditional.conf
# Only enable buildkitd if container engine is containerd.
[Service]
ExecCondition=test ${CONTAINER_ENGINE} = containerd
