# root/usr/local/lib/systemd/system/buildkitd.service.d/conditional.conf
# Only enable containerd if desired
[Service]
ExecCondition=test ${CONTAINER_ENGINE} = containerd
