# root/usr/local/lib/systemd/system/buildkitd.service.d/environment.conf
# Import Rancher Desktop configuration environment variables
[Service]
EnvironmentFile=-/etc/sysconfig/rancher-desktop
