# root/usr/local/lib/systemd/system/buildkitd.service.d/network-namespace.conf
[Unit]
Requires=network-namespace.service network-setup.service
After=network-setup.service

[Service]
ExecStart=
ExecStart=/usr/local/bin/wsl-exec \
  /usr/local/bin/buildkitd \
  --oci-worker=false --containerd-worker=true \
  --containerd-worker-addr=/run/k3s/containerd/containerd.sock \
  --addr fd://
